# odi.service
[Unit]
Description=Prometheus exporter for the HSGQ / ODI GPON SFP
Documentation=https://github.com/Strykar/GPON/
After=network-online.target
Wants=network-online.target
# Crashloop guard. Pairs with the exporter's behaviour of sys.exit-ing
# on AuthenticationException / BadHostKeyException: a wrong password
# causes the daemon to exit on the first cycle, systemd restarts it,
# the daemon exits again, and after 5 such restarts inside 60s systemd
# stops trying. Without this guard a wrong password would keep
# hammering the SFP (some firmwares lock the account out after a few
# failures and need a power-cycle to clear).
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple

# Password lives in /etc/gpon-exporter/credentials (mode 0600 root:root),
# one line: ONU_SSH_PASSWORD=your-password-here
# systemd parses EnvironmentFile= BEFORE switching to the dynamic UID,
# so the file can stay 0600 root-owned and the dynamic user never reads
# it directly.
EnvironmentFile=/etc/gpon-exporter/credentials

# --- Hardening ---------------------------------------------------------
# Run as a transient unprivileged user (UID assigned per start). The
# exporter doesn't need persistent on-disk identity beyond /var/lib/
# gpon-exporter (created automatically by StateDirectory= and owned by
# the dynamic UID with 0700).
DynamicUser=yes
StateDirectory=gpon-exporter
StateDirectoryMode=0700

# Strip everything we don't need. The exporter binds an unprivileged
# port, doesn't write outside its StateDirectory, doesn't talk to any
# kernel control interface, and only needs IPv4/IPv6 sockets.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
# If a future change extends the exporter to talk over a Unix socket
# (local IPC, /metrics over UDS instead of TCP), AF_INET/AF_INET6 only
# will return EAFNOSUPPORT inside the sandbox -- the failure points at
# code, not policy. Add AF_UNIX here when that happens. Same trap if
# we ever uncomment SystemCallFilter or MemoryDenyWriteExecute below
# and add bpf() / perf_event_open() / io_uring / JIT-using deps: the
# filter denies them by default. Check this unit file first when an
# inexplicable runtime error appears after a code change.
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
# MemoryDenyWriteExecute and SystemCallFilter omitted: the cryptography
# library that paramiko depends on uses runtime code generation in some
# configurations and a strict syscall allow-list often breaks paramiko's
# I/O paths. Enable these only after testing on your specific
# Python/paramiko/cryptography versions.

# --- ExecStart ---------------------------------------------------------
# /usr/local/bin is the FHS-correct location for sysadmin-installed
# scripts; adjust if you cloned the repo elsewhere. Webserver-port 8114
# is the documented default; change it if your prometheus.yml scrape
# config points elsewhere. --known-hosts lives under StateDirectory so
# the dynamic UID can write to it across restarts (StateDirectory
# preserves ownership).
ExecStart=/usr/bin/python3 /usr/local/bin/gpon_exporter.py \
    --device admin@192.168.1.1:22 \
    --webserver-port 8114 \
    --known-hosts ${STATE_DIRECTORY}/known_hosts

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
